# scripts/scanoss-api.service
[Unit]
Description=SCANOSS API Service

[Service]
ExecStart=/usr/bin/scanoss-api.sh
Restart=always